# kafka.conf
[program:kafka]
command=/kafka_2.10-0.8.2.1/bin/kafka-server-start.sh /kafka_2.10-0.8.2.1/config/server.properties
autorestart=true
priority=100
