# kafkadropwizard-supervisor.conf
[program:kafkarest]
directory=/opt/dropwizard-kafka
command=java -jar dropwizard-kafka.jar server /etc/dropwizard-kafka/kafka-http.yml
autorestart=true
priority=1000
