# zookeeper.conf
[program:zookeeper]
command=/kafka_2.10-0.8.2.1/bin/zookeeper-server-start.sh /kafka_2.10-0.8.2.1/config/zookeeper.properties
autorestart=true
priority=1
